# systemd/pi-tailscale-maintain.service
[Unit]
Description=Maintain Tailscale options, subnet routes, exit node, and performance settings
Wants=network-online.target
After=network-online.target tailscaled.service
Requires=tailscaled.service

[Service]
Type=oneshot
ExecStart=/root/pi-tailscale-role/bin/pi-tailscale-maintain.sh
