# systemd/pi-tailscale-health.service
[Unit]
Description=Lightweight Tailscale health check
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/root/pi-tailscale-role/bin/pi-tailscale-health.sh
